# scripts/com.da.overwatch.ai-patch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.da.overwatch.ai-patch</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/dykim/Documents/DAProj/(Clone) DA_Overwatch analysis/scripts/run_local_ai_patch_update.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/dykim/Documents/DAProj/(Clone) DA_Overwatch analysis</string>

  <key>RunAtLoad</key>
  <true/>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StartInterval</key>
  <integer>900</integer>

  <key>StandardOutPath</key>
  <string>/Users/dykim/Documents/DAProj/(Clone) DA_Overwatch analysis/logs/launchd_ai_patch.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/dykim/Documents/DAProj/(Clone) DA_Overwatch analysis/logs/launchd_ai_patch.err.log</string>
</dict>
</plist>
